# systemd/hermes-patch-guard.timer
[Unit]
Description=hermes-patch-manager drift guard timer

[Timer]
# Check shortly after boot and every 10 minutes thereafter.
OnBootSec=2min
OnUnitActiveSec=10min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
