# systemd/hermes-gateway.service.d/10-hermes-patch-manager.conf
# Installed by hermes-patch-manager install.sh as a drop-in for the Hermes
# gateway service. Injecting the loader via PYTHONPATH (outside the venv) is
# what makes patches survive `hermes update` rebuilding the venv.
#
# ExecStartPre uses a leading '-' so a patch problem is logged but NEVER blocks
# the gateway from starting.
[Service]
Environment=PYTHONPATH=/opt/hermes-patch-manager/loader
ExecStartPre=-/opt/hermes-patch-manager/hpm heal
